# build/files
/usr/local/bin/chronicle
/Library/LaunchDaemons/tk.unstac.chronicle.plist
